# crontab/crontab
# Yearly

# Monthly

# Weekly

# Daily

# Minutely
* * * * *  /home/emiel/dotfiles/crontab/batlow.sh
